# crontab
*/30 * * * * python /goes-notify.py >> /logfile